# docker/dockerfiles/dubbo-demo-consumer/supervisord.conf
[supervisord]
nodaemon=true

[program:demo-consumer]
command=/dubbo/demo-consumer/bin/start.sh