# docker/dockerfiles/dubbo-demo-provider/supervisord.conf
[supervisord]
nodaemon=true

[program:demo-provider]
command=/dubbo/demo-provider/bin/start.sh