# docker/dockerfiles/dubbo-monitor-simple/supervisord.conf
[supervisord]
nodaemon=true

[program:monitor-simple]
command=/dubbo/monitor-simple/bin/start.sh